# preconfiguredSettings/topBar/.config/autostart/volumeicon.desktop
[Desktop Entry]
Name=VolumeIcon
Exec=volumeicon
Encoding=UTF-8
Type=Application
Icon=volumeicon
