# preconfiguredSettings/launchers/autostart/autoTrash.desktop
[Desktop Entry]
Type=Application
Exec=trash-empty 40
TryExec=trash-empty 40
NoDisplay=true
