# preconfiguredSettings/launchers/autostart/numLockOn.desktop
[Desktop Entry]
Type=Application
Exec=xnumlock on
TryExec=xnumlock on
NoDisplay=true
