# preconfiguredSettings/launchers/autostart/redshift.desktop
[Desktop Entry]
Encoding=UTF-8
Version=1.0
Type=Application
Exec=bash -c "if [ -f /usr/bin/redshift ]; then while true; do if ! pgrep -xc redshift;then gtk-redshift -l $(cat /etc/latlon) || gtk-redshift; else for process in $(pgrep -f 'bash -c.*redshift');do kill -9 $process;break; done; killall redshift; fi; sleep 10; done; fi;"
Name=Redshift
Comment=Redshift
GenericName=Redshift
