# preconfiguredSettings/launchers/autostart/volumeIcon.desktop
[Desktop Entry]
Encoding=UTF-8
Version=1.0
Type=Application
Icon=system-run-symbolic
Name=VolumeIcon
GenericName=VolumeIcon
Comment=Volumeicon for alsa
Exec=bash -c "if [ -f /usr/bin/volumeicon ]; then while ! pgrep -xc volumeicon;do volumeicon;done; fi;"
NoDisplay=true
Terminal=false
