# preconfiguredSettings/launchers/autostart/xfce4-power-manager.desktop
[Desktop Entry]
Name=Power Manager
Icon=xfpm-ac-adapter
Exec=xfce4-power-manager
Terminal=false
Type=Application
StartupNotify=false
