# preconfiguredSettings/launchers/autostart/xscreensaver.desktop
[Desktop Entry]
Type=Application
Name=Screensaver
Icon=preferences-desktop-screensaver
Exec=xscreensaver -no-splash
TryExec=xscreensaver
